# ansible/roles/mastercontrol-config/files/mastercontrol-identity.service
[Unit]
Description=IrishEcho MasterControl Startup Identity
After=network.target
Before=openclaw.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mastercontrol-startup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target